# freepbx-17/install/etc/services.d/caddy/run
#!/command/execlineb -P
/usr/bin/caddy run --config /etc/caddy/Caddyfile
